# config/openclaw-gateway.service
[Unit]
Description=OpenClaw Gateway
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=__OPENCLAW_USER__
RuntimeDirectory=openclaw-gateway
RuntimeDirectoryMode=0700
LogsDirectory=openclaw
LogsDirectoryMode=0750
UMask=0077
Environment=HOME=/home/__OPENCLAW_USER__
Environment=NODE_ENV=production
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/local/bin/openclaw-gateway-launch
Restart=on-failure
RestartSec=30
RestartPreventExitStatus=78
TimeoutStartSec=90
TimeoutStopSec=30
SuccessExitStatus=0 143
OOMPolicy=continue
KillMode=control-group

[Install]
WantedBy=multi-user.target
