# config/openclaw-health.timer
[Unit]
Description=Periodic OpenClaw runtime health checks

[Timer]
OnBootSec=5m
OnCalendar=*:0/15
Persistent=true
AccuracySec=30s
RandomizedDelaySec=30s
Unit=openclaw-health.service

[Install]
WantedBy=timers.target
