# config/openclaw-backup.timer
[Unit]
Description=Daily verified OpenClaw backup

[Timer]
OnCalendar=*-*-* 03:15:00
Persistent=true
RandomizedDelaySec=15m
Unit=openclaw-backup.service

[Install]
WantedBy=timers.target
